# build/delugevpn.conf
[program:start-script]
autorestart = false
startsecs = 0
user = root
command = /root/start.sh
umask = 000

[program:watchdog-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/watchdog.sh
umask = 000
